# .config/systemd/user/brightness-sync.service
[Unit]
Description=Sync monitor brightness cache at login
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/brightness-queue init

[Install]
WantedBy=default.target
