# .config/systemd/user/lam-gui.service
[Unit]
Description=Arctis Manager tray
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/sh -lc 'sleep 2; exec /usr/bin/lam-gui --systray'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
